# autostart/qui-domains.desktop
[Desktop Entry]
Icon=qubes
Name=Domains Tray
Categories=System;Monitor;
Exec=python3 -mqui.tray.domains
Terminal=false
Type=Application