# autostart/qui-devices.desktop
[Desktop Entry]
Icon=media-removable
Name=Devices Tray
Categories=System;Monitor;
Exec=python3 -mqui.tray.devices
Terminal=false
Type=Application

